# overlay/system/etc/profile.d/devicehacks.sh
export QMLSCENE_DEVICE=haliumqsgcontext
