# overlay/system/usr/lib/systemd/user/snapd.aa-prompt-ui.service
[Unit]
Description=snapd UI prompt

[Service]
Type=dbus
BusName=io.snapcraft.Prompt
ExecStart=/usr/bin/snapd-aa-prompt-ui
